# 95-backlightd.rules
ACTION=="change", SUBSYSTEM=="drm", HOTPLUG=="1", RUN+="backlightctl refresh"
